# scripts/enter-gamemode.service
[Unit]
Description=Enter Gaming Mode

[Service]
Type=oneshot
ExecStart=/usr/local/bin/enter-gamemode.sh
